# 4363f978-5805-4ea7-bc07-55ff06a65466/var/lib/dpkg/info/freeradius.list
/.
/var
/var/log
/var/log/freeradius
/etc
/etc/init
/etc/init/freeradius.conf
/etc/init.d
/etc/init.d/freeradius
/etc/default
/etc/default/freeradius
/etc/pam.d
/etc/pam.d/radiusd
/etc/freeradius
/etc/freeradius/sites-enabled
/etc/freeradius/certs
/etc/freeradius/users
/etc/freeradius/clients.conf
/etc/freeradius/modules
/etc/freeradius/modules/smsotp
/etc/freeradius/modules/radutmp
/etc/freeradius/modules/redis
/etc/freeradius/modules/rediswho
/etc/freeradius/modules/ldap
/etc/freeradius/modules/realm
/etc/freeradius/modules/ippool
/etc/freeradius/modules/policy
/etc/freeradius/modules/wimax
/etc/freeradius/modules/sqlcounter_expire_on_login
/etc/freeradius/modules/expr
/etc/freeradius/modules/expiration
/etc/freeradius/modules/krb5
/etc/freeradius/modules/opendirectory
/etc/freeradius/modules/cache
/etc/freeradius/modules/sql_log
/etc/freeradius/modules/mac2vlan
/etc/freeradius/modules/soh
/etc/freeradius/modules/ntlm_auth
/etc/freeradius/modules/pap
/etc/freeradius/modules/always
/etc/freeradius/modules/dynamic_clients
/etc/freeradius/modules/unix
/etc/freeradius/modules/checkval
/etc/freeradius/modules/radrelay
/etc/freeradius/modules/logintime
/etc/freeradius/modules/attr_rewrite
/etc/freeradius/modules/digest
/etc/freeradius/modules/acct_unique
/etc/freeradius/modules/otp
/etc/freeradius/modules/files
/etc/freeradius/modules/detail
/etc/freeradius/modules/detail.example.com
/etc/freeradius/modules/preprocess
/etc/freeradius/modules/chap
/etc/freeradius/modules/detail.log
/etc/freeradius/modules/mschap
/etc/freeradius/modules/replicate
/etc/freeradius/modules/cui
/etc/freeradius/modules/exec
/etc/freeradius/modules/dhcp_sqlippool
/etc/freeradius/modules/sradutmp
/etc/freeradius/modules/linelog
/etc/freeradius/modules/pam
/etc/freeradius/modules/passwd
/etc/freeradius/modules/etc_group
/etc/freeradius/modules/smbpasswd
/etc/freeradius/modules/echo
/etc/freeradius/modules/mac2ip
/etc/freeradius/modules/perl
/etc/freeradius/modules/counter
/etc/freeradius/modules/attr_filter
/etc/freeradius/modules/inner-eap
/etc/freeradius/proxy.conf
/etc/freeradius/attrs.access_reject
/etc/freeradius/templates.conf
/etc/freeradius/policy.txt
/etc/freeradius/experimental.conf
/etc/freeradius/ldap.attrmap
/etc/freeradius/eap.conf
/etc/freeradius/sql.conf
/etc/freeradius/acct_users
/etc/freeradius/huntgroups
/etc/freeradius/attrs
/etc/freeradius/hints
/etc/freeradius/sqlippool.conf
/etc/freeradius/policy.conf
/etc/freeradius/attrs.pre-proxy
/etc/freeradius/preproxy_users
/etc/freeradius/attrs.access_challenge
/etc/freeradius/sites-available
/etc/freeradius/sites-available/default
/etc/freeradius/sites-available/README
/etc/freeradius/sites-available/inner-tunnel
/etc/freeradius/sites-available/proxy-inner-tunnel
/etc/freeradius/sites-available/dhcp.relay
/etc/freeradius/sites-available/control-socket
/etc/freeradius/sites-available/status
/etc/freeradius/sites-available/soh
/etc/freeradius/sites-available/vmps
/etc/freeradius/sites-available/virtual.example.com
/etc/freeradius/sites-available/robust-proxy-accounting
/etc/freeradius/sites-available/dhcp
/etc/freeradius/sites-available/dynamic-clients
/etc/freeradius/sites-available/decoupled-accounting
/etc/freeradius/sites-available/copy-acct-to-home-server
/etc/freeradius/sites-available/originate-coa
/etc/freeradius/sites-available/coa
/etc/freeradius/sites-available/buffered-sql
/etc/freeradius/sites-available/example
/etc/freeradius/attrs.accounting_response
/etc/logrotate.d
/etc/logrotate.d/freeradius
/usr
/usr/lib
/usr/lib/freeradius
/usr/lib/freeradius/rlm_eap_gtc-2.2.8.so
/usr/lib/freeradius/rlm_copy_packet-2.2.8.so
/usr/lib/freeradius/rlm_mschap-2.2.8.so
/usr/lib/freeradius/rlm_eap_peap-2.2.8.so
/usr/lib/freeradius/rlm_eap_md5-2.2.8.so
/usr/lib/freeradius/rlm_soh-2.2.8.so
/usr/lib/freeradius/rlm_attr_rewrite-2.2.8.so
/usr/lib/freeradius/rlm_linelog-2.2.8.so
/usr/lib/freeradius/rlm_radutmp-2.2.8.so
/usr/lib/freeradius/rlm_pap-2.2.8.so
/usr/lib/freeradius/rlm_cache-2.2.8.so
/usr/lib/freeradius/rlm_sql_log-2.2.8.so
/usr/lib/freeradius/rlm_logintime-2.2.8.so
/usr/lib/freeradius/rlm_eap-2.2.8.so
/usr/lib/freeradius/rlm_ippool-2.2.8.so
/usr/lib/freeradius/rlm_expr-2.2.8.so
/usr/lib/freeradius/rlm_realm-2.2.8.so
/usr/lib/freeradius/rlm_pam-2.2.8.so
/usr/lib/freeradius/rlm_dynamic_clients-2.2.8.so
/usr/lib/freeradius/rlm_checkval-2.2.8.so
/usr/lib/freeradius/rlm_sqlcounter-2.2.8.so
/usr/lib/freeradius/rlm_policy-2.2.8.so
/usr/lib/freeradius/rlm_counter-2.2.8.so
/usr/lib/freeradius/rlm_passwd-2.2.8.so
/usr/lib/freeradius/rlm_jradius-2.2.8.so
/usr/lib/freeradius/rlm_eap_mschapv2-2.2.8.so
/usr/lib/freeradius/rlm_eap_ttls-2.2.8.so
/usr/lib/freeradius/rlm_eap_tls-2.2.8.so
/usr/lib/freeradius/rlm_files-2.2.8.so
/usr/lib/freeradius/rlm_replicate-2.2.8.so
/usr/lib/freeradius/rlm_always-2.2.8.so
/usr/lib/freeradius/rlm_attr_filter-2.2.8.so
/usr/lib/freeradius/rlm_exec-2.2.8.so
/usr/lib/freeradius/rlm_fastusers-2.2.8.so
/usr/lib/freeradius/rlm_expiration-2.2.8.so
/usr/lib/freeradius/rlm_acctlog-2.2.8.so
/usr/lib/freeradius/rlm_unix-2.2.8.so
/usr/lib/freeradius/rlm_acct_unique-2.2.8.so
/usr/lib/freeradius/rlm_otp-2.2.8.so
/usr/lib/freeradius/rlm_digest-2.2.8.so
/usr/lib/freeradius/rlm_python-2.2.8.so
/usr/lib/freeradius/rlm_preprocess-2.2.8.so
/usr/lib/freeradius/rlm_dbm-2.2.8.so
/usr/lib/freeradius/rlm_chap-2.2.8.so
/usr/lib/freeradius/rlm_eap_sim-2.2.8.so
/usr/lib/freeradius/rlm_eap_leap-2.2.8.so
/usr/lib/freeradius/rlm_perl-2.2.8.so
/usr/lib/freeradius/rlm_wimax-2.2.8.so
/usr/lib/freeradius/rlm_detail-2.2.8.so
/usr/lib/freeradius/rlm_sqlippool-2.2.8.so
/usr/lib/freeradius/rlm_sql-2.2.8.so
/usr/share
/usr/share/doc
/usr/share/doc/freeradius
/usr/share/doc/freeradius/coding-methods.rst.gz
/usr/share/doc/freeradius/rlm_sqlcounter.gz
/usr/share/doc/freeradius/processing_users_file.rst
/usr/share/doc/freeradius/rlm_ldap.gz
/usr/share/doc/freeradius/snmp
/usr/share/doc/freeradius/release-method.rst
/usr/share/doc/freeradius/supervise-radiusd.rst.gz
/usr/share/doc/freeradius/rlm_eap.gz
/usr/share/doc/freeradius/performance-testing.gz
/usr/share/doc/freeradius/RADIUS-SQL.schema
/usr/share/doc/freeradius/rlm_python
/usr/share/doc/freeradius/Session-Type
/usr/share/doc/freeradius/bugs.gz
/usr/share/doc/freeradius/Post-Auth-Type
/usr/share/doc/freeradius/mssql
/usr/share/doc/freeradius/rlm_expiration
/usr/share/doc/freeradius/proxy.rst.gz
/usr/share/doc/freeradius/variables.rst.gz
/usr/share/doc/freeradius/MACOSX
/usr/share/doc/freeradius/rlm_sql.gz
/usr/share/doc/freeradius/rlm_sim_triplets
/usr/share/doc/freeradius/ascend
/usr/share/doc/freeradius/SoH.txt.gz
/usr/share/doc/freeradius/DIFFS.rst.gz
/usr/share/doc/freeradius/OS2
/usr/share/doc/freeradius/aaa.rst.gz
/usr/share/doc/freeradius/CYGWIN.rst.gz
/usr/share/doc/freeradius/copyright
/usr/share/doc/freeradius/README.gz
/usr/share/doc/freeradius/Autz-Type.rst
/usr/share/doc/freeradius/module_interface.rst.gz
/usr/share/doc/freeradius/tuning_guide
/usr/share/doc/freeradius/rlm_fastusers.gz
/usr/share/doc/freeradius/rlm_jradius
/usr/share/doc/freeradius/misc-nas.rst
/usr/share/doc/freeradius/bay
/usr/share/doc/freeradius/load-balance.rst.gz
/usr/share/doc/freeradius/README.rfc
/usr/share/doc/freeradius/Makefile.sphinx
/usr/share/doc/freeradius/examples
/usr/share/doc/freeradius/examples/clients.pl
/usr/share/doc/freeradius/examples/example.pl
/usr/share/doc/freeradius/examples/iplanet.schema
/usr/share/doc/freeradius/examples/openldap.schema
/usr/share/doc/freeradius/examples/users2mysql.pl
/usr/share/doc/freeradius/examples/radiusd.cron.monthly
/usr/share/doc/freeradius/examples/certs
/usr/share/doc/freeradius/examples/certs/client.cnf
/usr/share/doc/freeradius/examples/certs/README
/usr/share/doc/freeradius/examples/certs/server.cnf
/usr/share/doc/freeradius/examples/certs/bootstrap
/usr/share/doc/freeradius/examples/certs/xpextensions
/usr/share/doc/freeradius/examples/certs/ca.cnf
/usr/share/doc/freeradius/examples/certs/Makefile
/usr/share/doc/freeradius/examples/radwatch.in
/usr/share/doc/freeradius/examples/radiusd.cron.daily
/usr/share/doc/freeradius/examples/cryptpasswd
/usr/share/doc/freeradius/examples/radiusd.sh
/usr/share/doc/freeradius/examples/radiusd.cron.daily.in
/usr/share/doc/freeradius/examples/radsqlrelay
/usr/share/doc/freeradius/examples/min-includes.pl
/usr/share/doc/freeradius/examples/postgresql_update_radacct_group_trigger.sql
/usr/share/doc/freeradius/examples/cryptpasswd.in
/usr/share/doc/freeradius/examples/billing
/usr/share/doc/freeradius/examples/billing/README
/usr/share/doc/freeradius/examples/billing/clarent2db.pl
/usr/share/doc/freeradius/examples/billing/cisco_h323_db_sample_functions-postgres.sql
/usr/share/doc/freeradius/examples/billing/h323detail2db.pl
/usr/share/doc/freeradius/examples/billing/clarent_db_schema-postgres.sql
/usr/share/doc/freeradius/examples/radiusd2ldif.pl
/usr/share/doc/freeradius/examples/exec-program-wait
/usr/share/doc/freeradius/examples/radwatch
/usr/share/doc/freeradius/examples/create-users.pl
/usr/share/doc/freeradius/examples/rc.radiusd
/usr/share/doc/freeradius/examples/Makefile
/usr/share/doc/freeradius/examples/radiusd.cron.monthly.in
/usr/share/doc/freeradius/examples/iplanet.ldif
/usr/share/doc/freeradius/examples/rc.radiusd.in
/usr/share/doc/freeradius/ldap_howto.rst.gz
/usr/share/doc/freeradius/rlm_dbm.gz
/usr/share/doc/freeradius/RADIUS-LDAP-eDirectory
/usr/share/doc/freeradius/README.Debian
/usr/share/doc/freeradius/rlm_sqlippool
/usr/share/doc/freeradius/Simultaneous-Use.gz
/usr/share/doc/freeradius/configurable_failover.rst.gz
/usr/share/doc/freeradius/Acct-Type.rst
/usr/share/doc/freeradius/rlm_pam
/usr/share/doc/freeradius/rlm_krb5
/usr/share/doc/freeradius/CREDITS
/usr/share/doc/freeradius/rlm_passwd
/usr/share/doc/freeradius/index.rst
/usr/share/doc/freeradius/cisco.rst.gz
/usr/sbin
/usr/sbin/raddebug
/usr/sbin/freeradius
/usr/sbin/radwatch
/usr/sbin/checkrad
/usr/sbin/radmin
/usr/lib/freeradius/rlm_sqlippool.so
/usr/lib/freeradius/rlm_passwd.so
/usr/lib/freeradius/rlm_digest.so
/usr/lib/freeradius/rlm_expiration.so
/usr/lib/freeradius/rlm_logintime.so
/usr/lib/freeradius/rlm_attr_filter.so
/usr/lib/freeradius/rlm_python.so
/usr/lib/freeradius/rlm_dbm.so
/usr/lib/freeradius/rlm_files.so
/usr/lib/freeradius/rlm_linelog.so
/usr/lib/freeradius/rlm_fastusers.so
/usr/lib/freeradius/rlm_mschap.so
/usr/lib/freeradius/rlm_attr_rewrite.so
/usr/lib/freeradius/rlm_replicate.so
/usr/lib/freeradius/rlm_eap_md5.so
/usr/lib/freeradius/rlm_eap_peap.so
/usr/lib/freeradius/rlm_perl.so
/usr/lib/freeradius/rlm_otp.so
/usr/lib/freeradius/rlm_exec.so
/usr/lib/freeradius/rlm_ippool.so
/usr/lib/freeradius/rlm_acctlog.so
/usr/lib/freeradius/rlm_preprocess.so
/usr/lib/freeradius/rlm_counter.so
/usr/lib/freeradius/rlm_dynamic_clients.so
/usr/lib/freeradius/rlm_eap_sim.so
/usr/lib/freeradius/rlm_expr.so
/usr/lib/freeradius/rlm_eap.so
/usr/lib/freeradius/rlm_sqlcounter.so
/usr/lib/freeradius/rlm_copy_packet.so
/usr/lib/freeradius/rlm_sql_log.so
/usr/lib/freeradius/rlm_detail.so
/usr/lib/freeradius/rlm_acct_unique.so
/usr/lib/freeradius/rlm_eap_gtc.so
/usr/lib/freeradius/rlm_eap_mschapv2.so
/usr/lib/freeradius/rlm_policy.so
/usr/lib/freeradius/rlm_sql.so
/usr/lib/freeradius/rlm_soh.so
/usr/lib/freeradius/rlm_jradius.so
/usr/lib/freeradius/rlm_realm.so
/usr/lib/freeradius/rlm_checkval.so
/usr/lib/freeradius/rlm_radutmp.so
/usr/lib/freeradius/rlm_wimax.so
/usr/lib/freeradius/rlm_pam.so
/usr/lib/freeradius/rlm_pap.so
/usr/lib/freeradius/rlm_eap_tls.so
/usr/lib/freeradius/rlm_eap_leap.so
/usr/lib/freeradius/rlm_cache.so
/usr/lib/freeradius/rlm_chap.so
/usr/lib/freeradius/rlm_unix.so
/usr/lib/freeradius/rlm_always.so
/usr/lib/freeradius/rlm_eap_ttls.so
/usr/share/doc/freeradius/changelog.Debian.gz
