# 4363f978-5805-4ea7-bc07-55ff06a65466/var/lib/dpkg/info/freeradius.conffiles
/etc/default/freeradius
/etc/freeradius/acct_users
/etc/freeradius/attrs
/etc/freeradius/attrs.access_challenge
/etc/freeradius/attrs.access_reject
/etc/freeradius/attrs.accounting_response
/etc/freeradius/attrs.pre-proxy
/etc/freeradius/clients.conf
/etc/freeradius/eap.conf
/etc/freeradius/experimental.conf
/etc/freeradius/hints
/etc/freeradius/huntgroups
/etc/freeradius/ldap.attrmap
/etc/freeradius/modules/acct_unique
/etc/freeradius/modules/always
/etc/freeradius/modules/attr_filter
/etc/freeradius/modules/attr_rewrite
/etc/freeradius/modules/cache
/etc/freeradius/modules/chap
/etc/freeradius/modules/checkval
/etc/freeradius/modules/counter
/etc/freeradius/modules/cui
/etc/freeradius/modules/detail
/etc/freeradius/modules/detail.example.com
/etc/freeradius/modules/detail.log
/etc/freeradius/modules/dhcp_sqlippool
/etc/freeradius/modules/digest
/etc/freeradius/modules/dynamic_clients
/etc/freeradius/modules/echo
/etc/freeradius/modules/etc_group
/etc/freeradius/modules/exec
/etc/freeradius/modules/expiration
/etc/freeradius/modules/expr
/etc/freeradius/modules/files
/etc/freeradius/modules/inner-eap
/etc/freeradius/modules/ippool
/etc/freeradius/modules/krb5
/etc/freeradius/modules/ldap
/etc/freeradius/modules/linelog
/etc/freeradius/modules/logintime
/etc/freeradius/modules/mac2ip
/etc/freeradius/modules/mac2vlan
/etc/freeradius/modules/mschap
/etc/freeradius/modules/ntlm_auth
/etc/freeradius/modules/opendirectory
/etc/freeradius/modules/otp
/etc/freeradius/modules/pam
/etc/freeradius/modules/pap
/etc/freeradius/modules/passwd
/etc/freeradius/modules/perl
/etc/freeradius/modules/policy
/etc/freeradius/modules/preprocess
/etc/freeradius/modules/radrelay
/etc/freeradius/modules/radutmp
/etc/freeradius/modules/realm
/etc/freeradius/modules/redis
/etc/freeradius/modules/rediswho
/etc/freeradius/modules/replicate
/etc/freeradius/modules/smbpasswd
/etc/freeradius/modules/smsotp
/etc/freeradius/modules/soh
/etc/freeradius/modules/sql_log
/etc/freeradius/modules/sqlcounter_expire_on_login
/etc/freeradius/modules/sradutmp
/etc/freeradius/modules/unix
/etc/freeradius/modules/wimax
/etc/freeradius/policy.conf
/etc/freeradius/policy.txt
/etc/freeradius/preproxy_users
/etc/freeradius/proxy.conf
/etc/freeradius/sites-available/README
/etc/freeradius/sites-available/buffered-sql
/etc/freeradius/sites-available/coa
/etc/freeradius/sites-available/control-socket
/etc/freeradius/sites-available/copy-acct-to-home-server
/etc/freeradius/sites-available/decoupled-accounting
/etc/freeradius/sites-available/default
/etc/freeradius/sites-available/dhcp
/etc/freeradius/sites-available/dhcp.relay
/etc/freeradius/sites-available/dynamic-clients
/etc/freeradius/sites-available/example
/etc/freeradius/sites-available/inner-tunnel
/etc/freeradius/sites-available/originate-coa
/etc/freeradius/sites-available/proxy-inner-tunnel
/etc/freeradius/sites-available/robust-proxy-accounting
/etc/freeradius/sites-available/soh
/etc/freeradius/sites-available/status
/etc/freeradius/sites-available/virtual.example.com
/etc/freeradius/sites-available/vmps
/etc/freeradius/sql.conf
/etc/freeradius/sqlippool.conf
/etc/freeradius/templates.conf
/etc/freeradius/users
/etc/init.d/freeradius
/etc/init/freeradius.conf
/etc/logrotate.d/freeradius
/etc/pam.d/radiusd
