# 4363f978-5805-4ea7-bc07-55ff06a65466/var/lib/dpkg/info/hostapd.conffiles
/etc/default/hostapd
/etc/hostapd.conf
/etc/hostapd/ifupdown.sh
/etc/init.d/hostapd
