# 4363f978-5805-4ea7-bc07-55ff06a65466/var/lib/dpkg/info/linux-xenial-root-orangepizero.list
/.
/etc
/etc/armbian-release
/etc/X11
/etc/X11/xorg.conf.d
/etc/X11/xorg.conf.d/01-armbian-defaults.conf
/etc/fw_env.config
/etc/network
/etc/network/interfaces.hostapd
/etc/network/interfaces.default
/etc/network/interfaces.bonding
/etc/network/interfaces.network-manager
/etc/network/interfaces.r1switch
/etc/network/interfaces.r1
/etc/network/interfaces.espressobin
/etc/network/interfaces.r1router
/etc/kernel
/etc/kernel/preinst.d
/etc/kernel/preinst.d/initramfs-cleanup
/etc/profile.d
/etc/profile.d/ssh-title.sh
/etc/profile.d/activate_psd_user.sh
/etc/profile.d/check_first_login.sh
/etc/profile.d/check_first_login_reboot.sh
/etc/profile.d/90-vdpau.sh
/etc/init.d
/etc/init.d/firstrun
/etc/init.d/resize2fs
/etc/init.d/armhwinfo
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/71-no-recommends
/etc/apt/apt.conf.d/02periodic
/etc/cron.d
/etc/cron.d/armbian-updates
/etc/initramfs
/etc/initramfs/post-update.d
/etc/initramfs/post-update.d/99-uboot
/etc/default
/etc/default/armbian-motd.dpkg-dist
/etc/default/log2ram.dpkg-dist
/etc/default/cpufrequtils
/etc/modprobe.d
/etc/modprobe.d/r8723bs.conf
/etc/cron.daily
/etc/cron.daily/log2ram
/etc/lightdm
/etc/lightdm/lightdm.conf.d
/etc/lightdm/lightdm.conf.d/11-armbian.conf
/etc/lightdm/lightdm-gtk-greeter.conf
/etc/armbian.txt
/etc/update-motd.d
/etc/update-motd.d/98-autoreboot-warn
/etc/update-motd.d/41-armbian-config
/etc/update-motd.d/10-header
/etc/update-motd.d/35-tips
/etc/update-motd.d/40-updates
/etc/update-motd.d/30-sysinfo
/etc/update-motd.d/99-point-to-faq
/etc/mpv
/etc/mpv/mpv.conf
/usr
/usr/bin
/usr/bin/armbianmonitor
/usr/bin/sunxi_tp_temp
/usr/bin/h3consumption
/usr/bin/h3disp
/usr/share
/usr/share/armbian
/usr/share/armbian/boot.cmd
/usr/share/armbian/armbianEnv.txt
/usr/share/log2ram
/usr/share/log2ram/LICENSE
/usr/share/chromium
/usr/share/chromium/initial_bookmarks.html
/usr/share/chromium/master_preferences
/usr/lib
/usr/lib/NetworkManager
/usr/lib/NetworkManager/conf.d
/usr/lib/NetworkManager/conf.d/zz-override-wifi-powersave-off.conf
/usr/lib/nand-sata-install
/usr/lib/nand-sata-install/a10
/usr/lib/nand-sata-install/a10/bootloader
/usr/lib/nand-sata-install/a10/bootloader/u-boot.bin
/usr/lib/nand-sata-install/a10/bootloader/boot.ini
/usr/lib/nand-sata-install/a10/bootloader/linux
/usr/lib/nand-sata-install/a10/bootloader/linux/linux.ini
/usr/lib/nand-sata-install/a10/bootloader/boot.axf
/usr/lib/nand-sata-install/a13
/usr/lib/nand-sata-install/a13/bootloader
/usr/lib/nand-sata-install/a13/bootloader/boot.ini
/usr/lib/nand-sata-install/a13/bootloader/linux
/usr/lib/nand-sata-install/a13/bootloader/linux/u-boot.bin
/usr/lib/nand-sata-install/a13/bootloader/linux/linux.ini
/usr/lib/nand-sata-install/a13/bootloader/boot.axf
/usr/lib/nand-sata-install/a20
/usr/lib/nand-sata-install/a20/bootloader
/usr/lib/nand-sata-install/a20/bootloader/boot.ini
/usr/lib/nand-sata-install/a20/bootloader/linux
/usr/lib/nand-sata-install/a20/bootloader/linux/u-boot.bin
/usr/lib/nand-sata-install/a20/bootloader/linux/linux.ini
/usr/lib/nand-sata-install/a20/bootloader/boot.axf
/usr/lib/nand-sata-install/exclude.txt
/usr/lib/armbian
/usr/lib/armbian/apt-updates
/usr/lib/armbian/firstrun-config.sh
/usr/lib/chromium-browser
/usr/sbin
/usr/sbin/armbian-add-overlay
/usr/sbin/nand-sata-install
/usr/sbin/log2ram
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/71-axp-power-button.rules
/lib/systemd
/lib/systemd/system
/lib/systemd/system/getty@tty1.service.d
/lib/systemd/system/getty@tty1.service.d/10-noclear.conf
/lib/systemd/system/firstrun.service
/lib/systemd/system/haveged.service.d
/lib/systemd/system/haveged.service.d/10-no-new-privileges.conf
/lib/systemd/system/log2ram.service
/lib/systemd/system/firstrun-config.service
/lib/systemd/system/serial-getty@.service.d
/lib/systemd/system/serial-getty@.service.d/10-term.conf
/lib/systemd/system/systemd-modules-load.service.d
/lib/systemd/system/systemd-modules-load.service.d/10-timeout.conf
/lib/systemd/system/resize2fs.service
/boot
/boot/bin
/boot/bin/beelinkx2.bin
/boot/bin/lime2-emmc.bin
/boot/bin/pcduino3.bin
/boot/bin/lime2.bin
/boot/bin/orangepi.bin
/boot/bin/lamobo-r1.bin
/boot/bin/nanopim1plus.bin
/boot/bin/bananapim1pluslcd7.bin
/boot/bin/bananapilcd7.bin
/boot/bin/nanopiair.bin
/boot/bin/mk808c.bin
/boot/bin/bananapipro.bin
/boot/bin/lime.bin
/boot/bin/cubieboard2dual.bin
/boot/bin/nanopim1.bin
/boot/bin/aw-som-a20.bin
/boot/bin/bananapi.bin
/boot/bin/lime-a10-lcd.bin
/boot/bin/sunvell-r69.bin
/boot/bin/cubieboard.bin
/boot/bin/lime-a10.bin
/boot/bin/bananapim1plus.bin
/boot/bin/nanopineo.bin
/boot/bin/olinux-som-a13.bin
/boot/bin/olinuxino-a13.bin
/boot/bin/orangepione.bin
/boot/bin/micro.bin
/boot/bin/orangepizero.bin
/boot/bin/bananapim2plus.bin
/boot/bin/bananapiprolcd7.bin
/boot/bin/orangepilite.bin
/boot/bin/pcduino3nano.bin
/boot/bin/micro-emmc.bin
/boot/bin/orangepi-r1.bin
/boot/bin/orangepiplus2e.bin
/boot/bin/olimex-som-a20.bin
/boot/bin/cubietruck.bin
/boot/bin/orangepizeroplus2-h3.bin
/boot/bin/orangepipc.bin
/boot/bin/lime-a33.bin
/boot/bin/orangepipcplus.bin
/boot/bin/pcduino2.bin
/boot/bin/topwise-a721.bin
/boot/bin/orangepi2.bin
/boot/bin/cubieboard2.bin
/boot/bin/orangepiplus.bin
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/50-local.d
/var/lib/polkit-1/localauthority/50-local.d/printing.pkla
/var/lib/polkit-1/localauthority/50-local.d/power.pkla
/var/lib/polkit-1/localauthority/50-local.d/backlight.pkla
/var/lib/polkit-1/localauthority/50-local.d/plugdev.pkla
/var/lib/polkit-1/localauthority/50-local.d/networkmanager.pkla
/usr/lib/chromium-browser/master_preferences
