# 4363f978-5805-4ea7-bc07-55ff06a65466/var/lib/dpkg/info/sysstat.conffiles
/etc/cron.d/sysstat
/etc/cron.daily/sysstat
/etc/init.d/sysstat
/etc/sysstat/sysstat
/etc/sysstat/sysstat.ioconf
