# 4363f978-5805-4ea7-bc07-55ff06a65466/var/lib/dpkg/info/sysstat.list
/.
/var
/var/log
/var/log/sysstat
/etc
/etc/init.d
/etc/init.d/sysstat
/etc/cron.daily
/etc/cron.daily/sysstat
/etc/cron.d
/etc/cron.d/sysstat
/etc/sysstat
/etc/sysstat/sysstat.ioconf
/etc/sysstat/sysstat
/usr
/usr/bin
/usr/bin/iostat
/usr/bin/mpstat
/usr/bin/tapestat
/usr/bin/pidstat
/usr/bin/sar.sysstat
/usr/bin/sadf
/usr/bin/cifsiostat
/usr/share
/usr/share/doc
/usr/share/doc/sysstat
/usr/share/doc/sysstat/changelog.Debian.gz
/usr/share/doc/sysstat/README.gz
/usr/share/doc/sysstat/copyright
/usr/share/doc/sysstat/CREDITS.gz
/usr/share/doc/sysstat/README.Debian
/usr/share/doc/sysstat/xml
/usr/share/doc/sysstat/xml/sysstat.xsd.gz
/usr/share/doc/sysstat/xml/sysstat-3.1.dtd.gz
/usr/share/doc/sysstat/FAQ.gz
/usr/share/doc/sysstat/NEWS.Debian.gz
/usr/share/doc/sysstat/examples
/usr/share/doc/sysstat/examples/README.irqstat
/usr/share/doc/sysstat/examples/irqstat.gz
/usr/share/doc/sysstat/examples/crontab
/usr/share/doc/sysstat/examples/irqtop.gz
/usr/share/doc/sysstat/examples/sargraph2.gz
/usr/share/doc/sysstat/examples/sargraph.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pidstat.1.gz
/usr/share/man/man1/sar.sysstat.1.gz
/usr/share/man/man1/iostat.1.gz
/usr/share/man/man1/sadf.1.gz
/usr/share/man/man1/tapestat.1.gz
/usr/share/man/man1/cifsiostat.1.gz
/usr/share/man/man1/mpstat.1.gz
/usr/share/man/man5
/usr/share/man/man5/sysstat.5.gz
/usr/share/man/man8
/usr/share/man/man8/sadc.8.gz
/usr/share/man/man8/sa1.8.gz
/usr/share/man/man8/sa2.8.gz
/usr/share/doc-base
/usr/share/doc-base/sysstat-faq
/usr/lib
/usr/lib/sysstat
/usr/lib/sysstat/sa1
/usr/lib/sysstat/debian-sa1
/usr/lib/sysstat/sadc
/usr/lib/sysstat/sa2
